# crontab
#Backup Database
0 5 * * * /usr/bin/python3 /home/user/scripts/backup.py

#Weekly System Update
0 0 * * 0 /usr/bin/python3 /home/user/scripts/system_update.py

#Send Email Reports
30 8 * * 1-5 /usr/bin/python3 /home/user/scripts/email_reports.py

#Cleanup Temp Files
15 3 * * * /usr/bin/bash /home/user/scripts/cleanup_temp_files.sh

#Daily Database Sync
0 12 * * * /usr/bin/python3 /home/user/scripts/db_sync.py

#Monthly Backup
0 0 1 * * /usr/bin/python3 /home/user/scripts/monthly_backup.py

#Yearly Backup
0 0 1 1 * /usr/bin/python3 /home/user/scripts/yearly_backup.py
